# net-proxy/c-icap/files/c-icap
#!/sbin/openrc-run
# Copyright 2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="c-icap daemon"
description="c-icap-server"
command=/usr/bin/c-icap
command_args="${c-icap_args}"

depend() {
	need net
}

start_pre() {
	checkpath --directory --owner cicap:cicap --mode 0775 /var/run/c-icap
}

start() {
	ebegin "Starting c-icap-server"
	start-stop-daemon --start --exec $command --pidfile /var/run/c-icap/c-icap.pid
	eend $?

}

stop() {
	ebegin "Stopping c-icap-server"
	start-stop-daemon --stop --exec $command --pidfile /var/run/c-icap/c-icap.pid
	eend $?

}
